# crontab
0	*/6	*	*	*	${PWD}/get_games/get_games.sh
@reboot ¤${PWD}/BotPromo-Master/run.sh
