# crontab
0	*/6	*	*	*	${PWD}/display_games/run.sh
0	*/6	*	*	*	${PWD}/get_games/run.sh
0	0	1	*	*	${PWD}/clear_data.sh
