# crontab
0	*/6	*	*	*	$(pwd)/display_games/run.sh
0	*/6	*	*	*	$(pwd)/get_games/run.sh
0	0	1	*	*	$(pwd)/clear_data.sh
